# docs/install/fabric/configs/web/crontab
@reboot /opt/homelesshelper/daemons/start-all.sh